# supervisor.rsyslog.conf
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
stdout_logfile=syslog
stderr_logfile=syslog